# logger/loggserver.conf
description "Loggs Monitoring"
author "Edytarcio"

start on runlevel [2345]
stop on shutdown

respawn

script

    exec loggserver

end script

